# supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslog]
command=rsyslogd -n

[program:cron]
command=crond -f -S

[program:webserver]
command=python /ssl-provision/ssl.py
